# gpt-threads-compute/gpt-threads.service
[Unit]
After=network.service

[Service]
ExecStart=/home/user/gpt-threads/gpt-threads-compute/start_api.sh
StandardOutput=null
StandardError=null

[Install]
WantedBy=default.target