# gpt-threads-compute/task_cleanup.service
[Unit]
After=network.service

[Service]
ExecStart=/home/user/gpt-threads/gpt-threads-compute/auto_cleanup.sh
StandardOutput=null
StandardError=null

[Install]
WantedBy=default.target